# ansible/roles/jenkins-node/templates/jenkins-swarm.service.j2
[Unit]
Description=Jenkins swarm service
After=network.target

[Service]
Type=simple
User=iperf
ExecStart=/home/vagrant/jenkins-node.sh
Restart=always

[Install]
WantedBy=multi-user.target
